# assets/vps-proxy-manager/systemd/vps-proxy-manager.service
[Unit]
Description=VPS Proxy Manager Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vps-proxy-manager
EnvironmentFile=/etc/vps-proxy-manager/vps-proxy-manager.env
ExecStart=/opt/vps-proxy-manager/venv/bin/vps-proxy-manager run-bot
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vps-proxy-manager/data /etc/vps-proxy-manager

[Install]
WantedBy=multi-user.target
