# assets/vps-proxy-manager/systemd/vps-proxy-codex-worker.service
[Unit]
Description=VPS Proxy Manager Codex Provisioning Worker
After=network-online.target vps-proxy-manager.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vps-proxy-manager
EnvironmentFile=/etc/vps-proxy-manager/vps-proxy-manager.env
ExecStart=/opt/vps-proxy-manager/venv/bin/vps-proxy-manager run-codex-worker
Restart=on-failure
RestartSec=10
UMask=0077

[Install]
WantedBy=multi-user.target
